# packaging/DEBIAN/prerm
#!/bin/bash
set -e

case "$1" in
    remove|upgrade|deconfigure)
        systemctl stop PgTgPiGpio.service || true
        systemctl disable PgTgPiGpio.service || true
        ;;
esac
